# systemd/coredns-dnsredir@.service.in
#
# CoreDNS dnsredir systemd service file
#

[Unit]
Description=CoreDNS dnsredir __TAG_OR_COMMIT__
Documentation=https://github.com/leiless/dnsredir
Documentation=https://leiless.github.io/dnsredir/
Documentation=https://git.io/JJZ3N
After=network.target
Before=nss-lookup.target
Wants=nss-lookup.target

[Service]
Type=simple
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/local/bin/coredns-dnsredir -conf /etc/coredns-dnsredir/%i
LimitNOFILE=131072
Restart=on-failure
RestartSec=3
User=nobody
# For debugging
#StandardError=journal

[Install]
WantedBy=multi-user.target
